# timemachine/rootfs/etc/s6-overlay/s6-rc.d/smbd/finish
#!/command/with-contenv bashio
# ==============================================================================
# Cleanup on stop
# ==============================================================================
bashio::log.info "Samba stopped"
