# timemachine/rootfs/etc/s6-overlay/s6-rc.d/smbd/run
#!/command/with-contenv bashio
# ==============================================================================
# Start Samba daemon
# ==============================================================================
bashio::log.info "Starting Samba for Time Machine..."
exec smbd --foreground --debug-stdout --no-process-group </dev/null
