# timemachine/rootfs/etc/s6-overlay/s6-rc.d/avahi-daemon/run
#!/command/with-contenv bashio
bashio::log.info "Starting Avahi daemon (standalone, no D-Bus)..."

rm -f /var/run/avahi-daemon/pid

exec avahi-daemon --no-drop-root --no-chroot --no-rlimits --debug 2>&1
