# timemachine/rootfs/etc/s6-overlay/s6-rc.d/avahi-publish/run
#!/command/with-contenv bashio
declare share_name
share_name=$(bashio::config 'share_name')

bashio::log.info "Publishing Time Machine service via mDNS for share '${share_name}'..."

avahi-publish-service -f \
    "Home Assistant Time Machine" \
    _adisk._tcp 9 \
    "sys=waMa=0,adVF=0x100" \
    "dk0=adVN=${share_name},adVF=0x82" &

exec avahi-publish-service -f \
    "Home Assistant Time Machine" \
    _smb._tcp 445
